# rhasspy/profiles/en/supervisord.conf
[supervisord]
nodaemon=true

[program:microphone]
command=rhasspy-microphone-pyaudio-hermes --sample-rate 16000 --sample-width 2 --channels 1 --debug --host 10.0.0.76 --port 1883 --site-id mema01 --device-index 1 --udp-audio-host 127.0.0.1 --frames-per-buffer 1024
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:speakers]
command=rhasspy-speakers-cli-hermes --play-command 'aplay -q -t wav -D default:CARD=seeed2micvoicec' --list-command 'aplay -L' --volume 1 --debug --host 10.0.0.76 --port 1883 --site-id mema01
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:wake_word]
command=rhasspy-wake-porcupine-hermes --keyword porcupine_raspberry-pi.ppn --sensitivity 0.5 --keyword-dir /profiles/en/porcupine --debug --host 10.0.0.76 --port 1883 --site-id mema01
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:speech_to_text]
command=rhasspy-asr-kaldi-hermes --model-type nnet3 --model-dir /profiles/en/kaldi/model --graph-dir /profiles/en/kaldi/model/graph --debug --host 10.0.0.76 --port 1883 --site-id mema01 --dictionary /profiles/en/kaldi/dictionary.txt --language-model /profiles/en/kaldi/language_model.txt --language-model-type text_fst --base-dictionary /profiles/en/kaldi/base_dictionary.txt --base-dictionary /profiles/en/kaldi_custom_words.txt --dictionary-casing lower --g2p-model /profiles/en/kaldi/g2p.fst --unknown-words /profiles/en/kaldi_unknown_words.txt --base-language-model-fst /profiles/en/kaldi/base_language_model.fst --base-language-model-weight 0 --mixed-language-model-fst /profiles/en/kaldi/mixed.fst --frequent-words /usr/lib/rhasspy/rhasspy-profile/rhasspyprofile/profiles/en/frequent_words.txt --max-frequent-words 100 --max-unknown-words 8 --unknown-words-probability 1e-05 --unknown-token '<unk>' --silence-probability 0.5 --cancel-word '' --cancel-probability 0.01 --voice-skip-seconds 0 --voice-min-seconds 1 --voice-max-seconds 20 --voice-speech-seconds 0.3 --voice-silence-seconds 0.5 --voice-before-seconds 0.5 --voice-sensitivity 1 --voice-silence-method vad_only
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:intent_recognition]
command=rhasspy-nlu-hermes --intent-graph /profiles/en/intent_graph.pickle.gz --debug --host 10.0.0.76 --port 1883 --site-id mema01 --replace-numbers --language en_US --casing lower --converters-dir /profiles/en/converters --failure-token <unk>
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:text_to_speech]
command=rhasspy-tts-cli-hermes --tts-command 'bash -c '"'"'{%% if "/" in lang: %%}{%% set lang, voice = lang.split("/", maxsplit=1) %%}{%% endif %%} curl -sS -X GET -G --output - --data-urlencode INPUT_TYPE=TEXT --data-urlencode OUTPUT_TYPE=AUDIO --data-urlencode AUDIO=WAVE --data-urlencode LOCALE={{ lang }} {%% if voice: %%}--data-urlencode{%% endif %%} {%% if voice: %%}VOICE={{ voice }}{%% endif %%} --data-urlencode INPUT_TEXT="$0" http://10.0.0.76:59125/process --data-urlencode VOICE=en_UK/apope_low'"'"'' --voices-command 'curl -sS -X GET http://10.0.0.76:59125/voices' --language en-GB --use-jinja2 --volume 1 --debug --host 10.0.0.76 --port 1883 --site-id mema01
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:dialogue]
command=rhasspy-dialogue-hermes --debug --host 10.0.0.76 --port 1883 --site-id mema01 --session-timeout 30 --sound wake /usr/lib/rhasspy/etc/wav/beep_hi.wav --sound recorded /usr/lib/rhasspy/etc/wav/beep_lo.wav --sound error /usr/lib/rhasspy/etc/wav/beep_error.wav --volume 1 --min-asr-confidence 0 --say-chars-per-second 33 --sound-suffix .wav --sound-suffix .mp3 --sound-suffix .ogg --sound-suffix .flac
stopasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

